# packages/pegasus-ops/backup-agent/config/crontab
# Run every minute
*/1 * * * * /usr/local/bin/rclone_copy "localfile:/data/nfs/samba.backup/archon/Data Backup" gcs-secret:
