# packages/pegasus-ops/backup-agent/config/crontab
# Run every minute
*/1 * * * * /usr/local/bin/rclone_copy source:/data/archon/test dstFs=pegasus-backups:pegasus-backup-test/rclone_copy
