# packages/pegasus-ops/backup-agent/config/crontab
# Run every minute
*/1 * * * * echo "world"

# Run every 2 seconds
*/2 * * * * * * echo "hello"
*/2 * * * * * * echo "hello" 
